# systemd/color-detector.service
[Unit]
Description=Circle Detector - Color Detection Web Application
After=network.target

[Service]
Type=simple
User=sano
WorkingDirectory=/home/sano/color_detector_app
ExecStart=/usr/bin/python3 -m circle_detector.app
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
